# agent/yadro-snmp-agent.service.in
[Unit]
Description=Yadro SNMP Sub Agent
Wants=snmpd.service
After=snmpd.service
After=obmc-mapper.target

[Service]
Restart=always
Environment=OPTIONS="-Ls0-6d"
EnvironmentFile=-@sysconfdir@/default/yadro-snmp-agent
ExecStart=@bindir@/yadro-snmp-agent $OPTIONS
SyslogIdentifier=yadro-snmp

[Install]
WantedBy=@SYSTEMD_TARGET@
